# scripts/com.mlx-audio-transcriptor.watcher.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mlx-audio-transcriptor.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PYTHON__</string>
        <string>-m</string>
        <string>app.cli</string>
        <string>scan</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>

    <key>WatchPaths</key>
    <array>
        <string>__WATCH_DIR__</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/mlx-audio-transcriptor/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/mlx-audio-transcriptor/stderr.log</string>
</dict>
</plist>
